# etc/virtkvm.service
[Unit]
Description=Virtual KVM

[Service]
ExecStart=/usr/local/bin/virtkvm --daemon

[Install]
WantedBy=multi-user.target
